# home/zappyk/bin/deluge-web-custom-admin
#! /bin/bash
#
# deluge-web-custom-admin       Deluge Web Customer up/down
#
# chkconfig: - 01 99
# description: Activates/Deactivates Deluge Web Customer interfaces configured to \
#              start at boot time.
#
### BEGIN INIT INFO
# Provides: $network
# Should-Start: iptables ip6tables network
# Short-Description: Deluge Web Customer up/down
# Description: Deluge Web Customer up/down
### END INIT INFO

# Source function library.
. /etc/init.d/functions

rc=0

# See how we were called.
case "$1" in
  start)
	rc=0
	/root/bin/deluge-web-custom.sh | bash
	rc=$?
        ;;
  stop)
	rc=0
	/root/bin/deluge-web-custom-kill.sh $$ | bashh
	rc=$?
        ;;
  status)
	rc=0
	/root/bin/deluge-web-custom-kill.sh $$
	rc=$?
	;;
  restart|reload|force-reload)
	$0 stop
	$0 start
	rc=$?
	;;
  *)
        echo $"Usage: $0 {start|stop|status|restart|reload|force-reload}"
        exit 2
esac

exit $rc
