# home/zappyk/bin/raspberry-temp-check.cron
# Run the minutly jobs
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
* * * * * zappyk /home/bin/raspberry-temp-check.sh true
