# root/bin/zappyk-start
#!/bin/bash
#
# live: Init script for zappyk's autostart
#
# chkconfig: 345 99 00
# description: Init script for zappyk's autostart.

# . /etc/init.d/functions

# /root/bin/raspi1-mount-deluge-incomings

_stop____deluge() { /root/bin/deluge-web-custom-kill.sh; }
_start___deluge() { /root/bin/deluge-web-custom.sh     ; }
_status__deluge() { _stop____deluge; }
_restart_deluge() { _stop____deluge && _start___deluge; }

_stop____peervpn() { /root/bin/peervpn-1-client-kill
                     /root/bin/peervpn-2-server-kill; }
_start___peervpn() { #root/bin/peervpn-root-chroot.sh 2 server
                     /root/bin/peervpn-1-client
                     /root/bin/peervpn-2-server; }
_status__peervpn() { _stop____peervpn; }
_restart_peervpn() { _stop____peervpn && _start___peervpn; }
_tail____peervpn() { echo '#==========================================================='
                     echo '#   tail /var/log/peervpn-1-client.log                      '
                     echo '#___________________________________________________________'
                     echo '#                                                           '
                     tail /var/log/peervpn-1-client.log
                     echo '#==========================================================='
                     echo '#   tail /var/log/peervpn-2-server.log                      '
                     echo '#___________________________________________________________'
                     echo '#                                                           '
                     tail /var/log/peervpn-2-server.log; }

case "$1" in
    stop            ) $0 stop-deluge
                      $0 stop-peervpn
                    ;;
    start           ) $0 start-deluge
                    # $0 start-peervpn
                    ;;
    status          ) $0 status-deluge
                      $0 status-peervpn
                    ;;
    restart         ) $0 restart-deluge
                    # $0 restart-peervpn
                    ;;
    tail            ) $0 tail-peervpn
                    ;;
    stop-deluge     ) _stop____deluge  | bash ;;
    stop-peervpn    ) _stop____peervpn | bash ;;
    start-deluge    ) _start___deluge  | bash ;;
    start-peervpn   ) _start___peervpn ;;
    status-deluge   ) _status__deluge  ;;
    status-peervpn  ) _status__peervpn ;;
    restart-deluge  ) _restart_deluge  ;;
    restert-peervpn ) _restart_peervpn ;;
    tail-peervpn    ) _tail____peervpn ;;
    *               ) echo -e "Used:\n\t$0 [ stop | start | status | restart | tail ] (-deluge|-peervpn)" && exit 1 ;;
esac

exit
